# autostart/nvfancontrol.desktop
[Desktop Entry]
Exec=nvfancontrol
Icon=
Name=nvfancontrol
Path=
Terminal=False
Type=Application
